# workbuddy-side/launchd/com.workbuddy.hub-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workbuddy.hub-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>WORKBUDDY_SKILL_DIR/workbuddy-side/scripts/hub_watcher.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HUB_URL</key>
        <string>http://localhost:3100</string>
        <key>HUB_AGENT_ID</key>
        <string>workbuddy</string>
        <key>HUB_WATCHER_LOG</key>
        <string>INFO</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>WORKBUDDY_SKILL_DIR</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/hub-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/hub-watcher.err</string>
</dict>
</plist>
